# php-modern/s6/s6-overlay/s6-rc.d/php81/run
#!/command/execlineb -P

/usr/local/php81/sbin/php-fpm