# php-modern/s6/s6-overlay/s6-rc.d/php73/run
#!/command/execlineb -P

/usr/local/php73/sbin/php-fpm