# php-modern/s6-overlay/services.d/php81/run
#!/usr/bin/with-contenv sh

exec /usr/local/php81/sbin/php-fpm