# s6-overlay/services.d/php72/run
#!/usr/bin/with-contenv sh

exec /usr/local/php72/sbin/php-fpm