# php-modern/s6/s6-overlay/s6-rc.d/php70/run
#!/command/execlineb -P

/usr/local/php70/sbin/php-fpm
