# s6-overlay/services.d/crond/run
#!/usr/bin/execlineb -P
/usr/sbin/crond -f -L /usr/local/var/log/cron/crond.log