# openresty125/s6/s6-overlay/s6-rc.d/nginx/run
#!/bin/sh

/usr/local/openresty/nginx/sbin/nginx -g 'daemon off;'