# s6-overlay/services.d/nginx/run
#!/usr/bin/execlineb -P
/usr/local/openresty/nginx/sbin/nginx -g 'daemon off;'