# php-modern/s6-overlay/services.d/php70/run
#!/usr/bin/with-contenv sh

exec /usr/local/php70/sbin/php-fpm
